# infra/ngrok.service
[Unit]
Description=ngrok Tunnel for SavageGPT
After=savagegpt.service

[Service]
User=deployer
WorkingDirectory=/opt/savagegpt
EnvironmentFile=/etc/savagegpt.env
ExecStart=/usr/local/bin/ngrok http --authtoken ${NGROK_AUTH_TOKEN} --name savagegpt 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
