# infra/ngrok.service
[Unit]
Description=ngrok Tunnel for SavageGPT
After=savagegpt.service

[Service]
User=asehran
WorkingDirectory=/opt/savagegpt
EnvironmentFile=/etc/savagegpt.env
ExecStart=/usr/local/bin/ngrok http 8000 \
  --authtoken ${NGROK_AUTH_TOKEN} \
  --name savagegpt \
  --log=stdout
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
