# infra/savagegpt.service
[Unit]
Description=SavageGPT FastAPI Service
After=network.target

[Service]
User=deployer
WorkingDirectory=/opt/savagegpt
EnvironmentFile=/etc/savagegpt.env
ExecStart=/usr/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
